# meta-megabit/recipes-apps/blink/files/blink.service
[Unit]
Description=Blink startup script

[Service]
ExecStart=/usr/bin/pi-blink

[Install]
WantedBy=multi-user.target